# rootfs/etc/supervisor.d/haproxy.ini
[program:haproxy]
command=/usr/local/sbin/haproxy -f /usr/local/etc/haproxy/haproxy.cfg
numprocs=1
autostart=true
autorestart=true
